# stats/stats.cron
0 1 * * * /usr/local/bin/python3 /usr/src/app/fetcher/run.py -c 'https://raw.githubusercontent.com/CSSEGISandData/COVID-19/master/csse_covid_19_data/csse_covid_19_time_series/time_series_covid19_confirmed_global.csv' -e 'https://raw.githubusercontent.com/CSSEGISandData/COVID-19/master/csse_covid_19_data/csse_covid_19_time_series/time_series_covid19_deaths_global.csv' >> /var/log/cron.log 2>&1
